# roles/server-nas-slave/templates/sickrage.conf.j2
[program:sickrage]
priority=12
directory={{ sickrage_app_path }}
command=python SickBeard.py -p 8081 --nolaunch --datadir={{ sickrage_config_path }}
user={{ user }}
